# harness/systemd/terrarium-chronicler.timer
[Unit]
Description=Review the terrarium's evolution once an hour

[Timer]
OnActiveSec=12min
OnUnitInactiveSec=1h
AccuracySec=2min
Persistent=false

[Install]
WantedBy=timers.target
